# boot/animation.service
[Unit]
Description=Display Boot Animation
DefaultDependencies=no
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/path/to/your/script.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target